# xserver/app/conf.d/xvnc.conf
[program:xvnc]
command=Xvnc -ac -pn -dpms -depth 24 %(ENV_DISPLAY)s -rfbport=5901 -SecurityTypes=None -desktop=RobotOps
autorestart=true
stdout_logfile=/app/xvnc.log
redirect_stderr=true