# xserver/conf.d/xsession.conf
[program:xsession]
command=/etc/X11/Xsession
autorestart=true
stdout_logfile=/opt/xsession.log
redirect_stderr=true